# deploy/gardenops.service.example
[Unit]
Description=GardenOps API
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
User=gardenops
Group=gardenops
WorkingDirectory=/srv/gardenops/current
EnvironmentFile=/etc/gardenops.env
ExecStartPre=/usr/local/sbin/gardenops-release-preflight /srv/gardenops/current
ExecStart=/srv/gardenops/current/.venv/bin/python -m uvicorn gardenops.main:app --host 127.0.0.1 --port 8000 --no-proxy-headers
Restart=always
RestartSec=5
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/opt/gardenops/logs /opt/gardenops/media_uploads

[Install]
WantedBy=multi-user.target
